# configs/gunicorn_restoran.socket
[Unit]
Description=gunicorn socket for restoran project

[Socket]
ListenStream=/run/gunicorn/restoran_socket

[Install]
WantedBy=sockets.target
